# supervisor.conf
[supervisord]
nodaemon=true

# Start ollama server
[program:ollama]
command=ollama serve

# Set the command to run the Flask app
[program:main_server]
command=python ./src/main_server.py

